# init/savehouse.conf
# this is the startup for savehouse.py, it catches
# messages from mqtt and saves device data to my database
# optional stuff
description "desert-home iriscontrol python script"
version "1.0"
author "Dave"

start on (local-filesystems and net-device-up IFACE=eth0)
respawn
stop on shutdown
# configuration variables.
# haven't got any yet

script
# My startup script, plain old shell scripting here.
exec /home/pi/src/house/savehouse.py >> /var/log/house/savehouse.log 2>&1
end script
